# rootfs/etc/services.d/ps3netsrv/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid abc ps3netsrv /games
