# rootfs/etc/services.d/ps3netsrv/run
#!/usr/bin/execlineb -P
s6-setuidgid ps3netsrv
/usr/local/bin/ps3netsrv /games
